# sysadmin_scripts/ec2servers/copy-out-files/db/crontab
24 3 * * * sh /root/backup/run-ec2-consistent-snapshot-master.sh >> /root/backup/logs/log 2>&1
